# crontabs/root
# Reboot at 4:30am every day
# Note: To avoid infinite reboot loop, wait 70 seconds
# and touch a file in /etc so clock will be set
# properly to 4:31 on reboot before cron starts.
30 3 * * * eval "echo root; sleep 3; echo root; sleep 3; echo reboot; sleep 5;" | telnet 192.168.8.1
40 3 * * *  eval "echo root; sleep 3; echo root; echo 'ip link set dev ep0 qlen 50000'; echo 'ip link set dev ep0  mtu 2000'; echo 'ip link set dev ep0  up'; echo 'ip link set dev pon0 qlen 50000'; echo 'ip link set dev pon0  mtu 2000'; echo 'ip link set dev pon0 up'; echo 'ip link set dev eth0 qlen 50000'; echo 'ip link set dev eth0  mtu 2000'; echo 'ip link set dev eth0 up'; echo 'ip link set dev br0 qlen 50000'; echo 'ip link set dev br0 up'; echo 'ip link set dev br0  mtu 2000'; echo 'ip link set dev eth0.0 qlen 50000';echo 'ip link set dev eth0.0  mtu 2000'; echo 'ip link set dev eth0.0 up'; echo 'ip link set dev pon0.0 qlen 50000';echo 'ip link set dev pon0.0  mtu 2000'; echo 'ip link set dev pon0.0 up'; sleep 5;" | telnet 192.168.8.1 
05 3 * * * sleep 70 && touch /etc/banner && reboot
