# crontabs/root
# Reboot at 4:30am every day
# Note: To avoid infinite reboot loop, wait 70 seconds
# and touch a file in /etc so clock will be set
# properly to 4:31 on reboot before cron starts.
# 05 3 * * * eval "echo admin; sleep 1; echo admin123; sleep 1; echo reboot; sleep 5;" | telnet 192.168.2.1
05 3 * * * sleep 70 && touch /etc/banner && reboot
*/5 * * * * /bin/sh /etc/verify_ipv6.sh
00 2 * * * curl --output /etc/unbound/root.hints https://www.internic.net/domain/named.cache
