# rc/rc.local
# Put your custom commands here that should be executed once \n
# the system init finished. By default this file does nothing.

echo '1600' > /sys/kernel/slab/fq_flow_cache/cpu_partial
echo 400 > /sys/kernel/slab/TCP/cpu_partial
echo 400 > /sys/kernel/slab/TCPv6/cpu_partial
echo 400 > /sys/kernel/slab/UDP/cpu_partial
echo 400 > /sys/kernel/slab/UDPv6/cpu_partial
echo 400 > /sys/kernel/slab/UDPLITEv6/cpu_partial
echo 400 > /sys/kernel/slab/PING/cpu_partial
echo 400 > /sys/kernel/slab/PINGv6/cpu_partial
echo 400 > /sys/kernel/slab/wg_peer/cpu_partial
echo "nameserver 172.64.36.1" > /etc/resolv.conf

exit 0
