# rc/rc.local
# Put your custom commands here that should be executed once \n
# the system init finished. By default this file does nothing.
/etc/jool/setupjool.sh
sleep 30
echo performance > /sys/devices/system/cpu/cpufreq/policy0/scaling_governor
tc qdisc add dev eth0 root noqueue
tc qdisc add dev br-lan root noqueue
tc qdisc add dev eth1 root noqueue
#ethtool -A eth1 rx off tx off
#ethtool -A eth0 rx off tx off

exit 0
